# mkosi.extra/etc/systemd/system/journal-bridge.service
[Unit]
Description=Journal to File Bridge for remote_syslog
Documentation=https://github.com/tinkerbell/hook
After=systemd-journald.service remote-syslog.service
Requires=remote-syslog.service

[Service]
Type=simple
# Stream tink-agent journal entries to a file for remote_syslog to tail
# Use -o cat for raw messages (no timestamp prefix), stdbuf for line buffering
ExecStartPre=/bin/mkdir -p /run/log
ExecStart=/bin/sh -c 'exec stdbuf -oL journalctl -f -n 1000 -o cat -p info \
    -u tink-agent.service -u tink-agent-setup.service \
    >> /run/log/journal.log'
Restart=always
RestartSec=5

# Cleanup on stop
ExecStopPost=/bin/rm -f /run/log/journal.log

[Install]
WantedBy=multi-user.target
