# mkosi.extra/etc/systemd/system/remote-syslog.service
[Unit]
Description=Remote Syslog Log Forwarder
Documentation=https://github.com/papertrail/remote_syslog2
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/remote-syslog-start
Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/run/log /run

[Install]
WantedBy=multi-user.target
