# mkosi.extra/etc/systemd/system/tink-agent-setup.service
[Unit]
Description=Tink Agent Setup (pull image and extract binary)
After=containerd.service network-online.target systemd-resolved.service
Wants=network-online.target
Requires=containerd.service
# Only start when Tinkerbell parameters are on the kernel command line
ConditionKernelCommandLine=|tink_worker_image
ConditionKernelCommandLine=|docker_registry

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/tink-agent-setup
TimeoutStartSec=300

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tink-agent-setup

[Install]
WantedBy=multi-user.target
